# task1/80-test.rules
SUBSYSTEM=="usb", ACTION == "ADD", ENV{DEVTYPE}=="usb_device", RUN+="/bin/copylog.sh"
