# task1/80-test.rules
ACTION=="add", ATTRS{serial}=="408D5CBEC94AB471998546ED",RUN+="/usr/bin/add_trigger.sh"
